# inferno/config/inferno.conf
description     "inferno"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/inferno -s discomaster01
